# util-scripts/desu-timezone/desu-timezone.timer
[Unit]
Description=15min Timer

[Timer]
OnBootSec=0min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target